# config/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/nerv/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info
childlogdir=/var/log/nerv

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=5
startsecs=3
stderr_logfile=/var/log/nerv/nginx_error.log
stdout_logfile=/var/log/nerv/nginx_access.log
stderr_logfile_maxbytes=50MB
stdout_logfile_maxbytes=50MB
priority=100

[program:nerv-ui]
command=python3 -m http.server 3000 --bind 0.0.0.0 --directory /opt/nerv-ui
autostart=true
autorestart=true
startretries=5
startsecs=2
stderr_logfile=/var/log/nerv/ui_stderr.log
stdout_logfile=/var/log/nerv/ui_stdout.log
stderr_logfile_maxbytes=20MB
stdout_logfile_maxbytes=20MB
priority=200

[program:comfyui]
command=python3 %(ENV_COMFYUI_DIR)s/main.py --listen 0.0.0.0 --port %(ENV_COMFYUI_PORT)s --preview-method auto --fp16-vae
directory=%(ENV_COMFYUI_DIR)s
autostart=true
autorestart=unexpected
startretries=3
startsecs=10
exitcodes=0
stopwaitsecs=30
stderr_logfile=/var/log/nerv/comfyui_stderr.log
stdout_logfile=/var/log/nerv/comfyui_stdout.log
stderr_logfile_maxbytes=100MB
stdout_logfile_maxbytes=100MB
priority=300
environment=
    PYTHONUNBUFFERED="1",
    CUDA_VISIBLE_DEVICES="%(ENV_CUDA_VISIBLE_DEVICES)s"

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
